# linux/quetoo/usr/lib/systemd/system/quetoo-dedicated@.service
[Unit]
Description=Quetoo Dedicated Server (%i)
After=network.target
Documentation=https://github.com/jdolan/quetoo

[Service]
Type=forking
User=quetoo
Group=quetoo
# Runs inside a named screen session so admins can attach:
#   screen -r quetoo-%i
ExecStart=/usr/bin/screen -dmS quetoo-%i /usr/lib/quetoo/bin/quetoo-dedicated-start %i
ExecStop=/usr/bin/screen -S quetoo-%i -X quit
KillMode=control-group
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
